# scripts/cron.d/pbx3-recordings.example
# Call recordings — local archive offload + retention (R1.5) + S3 upload (S7).
# Install: cp to /etc/cron.d/pbx3-recordings and adjust user/paths.
#
# Replaces legacy pbx3 shell jobs (offload_simple.sh, agerecordings.sh,
# agegracerecordings.sh, manageRecs.php) — remove those from /etc/cron.d/pbx3.
#
# Env (optional) in /opt/pbx3api/.env:
#   PBX3_RECORDINGS_OFFLOAD_ENABLED (default true)
#   PBX3_RECORDINGS_STABLE_MINUTES  (default 5)
#   PBX3_RECORDINGS_ARCHIVE_ROOT    (default /opt/pbx3/media/recordings)
#   PBX3_RECORDING_UPLOAD_ENABLED   (default false — set true for S7)
#   PBX3_GATEKEEPER_URL             (e.g. https://control.pbx3.com)
#   PBX3_GATEKEEPER_TOKEN           (Bearer accepted by gatekeeper)
#   PBX3_RECORDING_UPLOAD_TENANTS   (optional allowlist, e.g. duns)

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Offload stable spool wavs to the local archive and index rows (every 10 min).
*/10 * * * * root cd /opt/pbx3api && /usr/bin/php artisan pbx3:recordings-offload >> /var/log/pbx3-recordings.log 2>&1

# S7 — upload indexed local archive rows to dedicated recordings bucket (same cadence).
*/10 * * * * root cd /opt/pbx3api && /usr/bin/php artisan pbx3:recordings-s3-upload >> /var/log/pbx3-recordings.log 2>&1

# Retention: age-out + grace purge + recmaxsize + recused tally (daily 02:30).
30 2 * * * root cd /opt/pbx3api && /usr/bin/php artisan pbx3:recordings-retain >> /var/log/pbx3-recordings.log 2>&1
